# meta-ampere/meta-common/recipes-ampere/nic/card-manager/xyz.openbmc_project.card.manager.service.in
[Unit]
Description= NIC management
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart=@bindir@/card-manager
Type=dbus
BusName=xyz.openbmc_project.card.manager
SyslogIdentifier=card-manager

[Install]
WantedBy=multi-user.target
