# meta-ampere/meta-jade/recipes-phosphor/gpio/host-gpio-handling/host-gpio-handling.service
[Unit]
Description=Dummy service

[Service]
ExecStart=/bin/sh -c "echo DoNothing"