# meta-ampere/meta-mitchell/recipes-ampere/host/ampere-sysfw-hang-handler/ampere-sysfw-hang-handler.service
[Unit]
Description=Ampere System Firmware Hang Handler
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/env ampere_sysfw_hang_handler.sh

[Install]
WantedBy=multi-user.target
